# root/etc/init.d/pool
#!/bin/sh /etc/rc.common

START=99
STOP=10

start() {
    ENABLED=$(uci get pool.setting.enable_feature 2>/dev/null)
    [ "$ENABLED" = "1" ] || return 0

    echo "Starting VLESS proxy..."
    # 这里启动 xray 或其他代理程序（示例）
    /usr/bin/xray -config /etc/xray/config.json &
}

stop() {
    echo "Stopping VLESS proxy..."
    killall xray 2>/dev/null
}
